# upstart/apu_lcm_tunnel_server.conf
description "bot-lcm-tunnel-server"

# Command
exec sudo -u husky bash -c 'source ~/.bashrc; sleep 2; bot-lcm-tunnel 192.168.1.11'

start on runlevel [2345]
stop on runlevel [016]
respawn
